# backend/deploy/supervisord.production.conf
; Solomon AI — production supervisor config for the backend
;
; BLOCKER #7 from production audit. Replaces the dev-mode preview pod
; config (`uvicorn server:app --workers 1 --reload`) which is unsafe for
; real Sunday-morning load.
;
; Drop this file at /etc/supervisor/conf.d/backend-prod.conf on whatever
; host runs the production backend (NOT in the Emergent preview pod —
; the platform manages that file). Reload supervisor:
;     sudo supervisorctl reread
;     sudo supervisorctl update
;     sudo supervisorctl restart backend
;
; Required env vars at the host level:
;   - ENVIRONMENT=production
;   - MONGO_URL, DB_NAME, STRIPE_API_KEY, STRIPE_PUBLISHABLE_KEY,
;     STRIPE_WEBHOOK_SECRET (whsec_...), EMERGENT_LLM_KEY, RESEND_API_KEY,
;     SOLOMON_SEED_PASSWORD (random, 32+ chars; not "change_me")

[program:backend]
command=/root/.venv/bin/gunicorn server:app -c /app/backend/gunicorn.conf.py
directory=/app/backend
autostart=true
autorestart=true
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true
environment=ENVIRONMENT="production",GUNICORN_WORKERS="4",GUNICORN_TIMEOUT="60"
stderr_logfile=/var/log/supervisor/backend.err.log
stdout_logfile=/var/log/supervisor/backend.out.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
